# ansible/files/configs/udev/95-monitors.rules
ACTION=="change", RUN+="/home/luiz/.monitor-hot-plug"
